# dist/etc/systemd/system/mpi-gunicorn.service
[Unit]
Description=Start MPI action web services
After=network.target

[Service]
Type=simple
Environment=LC_ALL=C.UTF-8
Environment=LANG=C.UTF-8
WorkingDirectory=/usr/lib/mpi
ExecStart=/usr/lib/mpi/bin/gunicorn

[Install]
WantedBy=multi-user.target