# dist/etc/systemd/system/mpi-process.service
[Unit]
Description=Search for new auctions and notify any subscriptions
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/usr/lib/mpi
ExecStart=/usr/lib/mpi/bin/console process

[Install]
WantedBy=multi-user.target