# dist/etc/systemd/system/mpi-gunicorn.service
[Unit]
Description=Start MPI Auction Web
After=network.target

[Service]
Type=simple
ExecStart=/usr/lib/mpi/bin/gunicorn

[Install]
WantedBy=multi-user.target